# Library/LaunchAgents/jp.mml.brewfileupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>jp.mml.brewfileupdate.plist</string>
	<key>Program</key>
	<string>/Users/kata/bin/brewcaskupdate.sh</string>
	<key>WorkingDirectory</key>
	<string>/Users/kata/tmp</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>00</integer>
		<key>Minute</key>
		<integer>05</integer>
	</dict>
</dict>
</plist>
